# templates/ram_disk.service.erb
[Unit]
Description=RAM disk

[Service]
Type=oneshot
ExecStart=mount -t tmpfs -o size=<%= ram_size %> tmpfs <%= record_dir %>

[Install]
WantedBy=multi-user.target